# roles/network/files/docker-gen-dnsmasq.conf
start on started docker
stop on stopped docker

respawn

exec docker-gen -watch=true -notify="/etc/init.d/dnsmasq force-reload" /etc/templates/dockerhosts.tmpl /etc/dockerhosts

