# roles/network/files/docker-gen-nginx.conf
start on started docker
stop on stopped docker

respawn

exec docker-gen -watch=true -notify="/etc/init.d/nginx reload" /etc/templates/nginx.tmpl /etc/nginx/sites-enabled/default

